# os/package/drafter-liveness/drafter-liveness.service.in
[Unit]
Description=Drafter liveness daemon
After=network.target @DRAFTER_LIVENESS_SYSTEMD_DEPENDENCY@
Requires=@DRAFTER_LIVENESS_SYSTEMD_DEPENDENCY@

[Service]
Type=oneshot
ExecStart=/usr/bin/drafter-liveness --vsock-port @DRAFTER_LIVENESS_VSOCK_PORT@
StandardOutput=file:/dev/kmsg
StandardError=file:/dev/kmsg

[Install]
WantedBy=multi-user.target
